# etc/local.d/zram.stop
/sbin/swapoff /dev/zram0; \
/bin/echo 1 > /sys/block/zram0/reset
